# contrib/pwsw.service
[Unit]
Description=PWSW - PipeWire Switcher daemon
Documentation=https://github.com/ledati16/pwsw
After=pipewire.service
Requires=pipewire.service
BindsTo=pipewire.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/bin/pwsw daemon --foreground
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pwsw

[Install]
WantedBy=graphical-session.target
